# lib/systemd/system/studio-stats.timer
[Timer]
OnCalendar=hourly
RandomizedDelaySec=300
AccuracySec=1us

[Install]
WantedBy=timer.target
