# lib/systemd/system/studio-clearsessions.timer
[Timer]
OnCalendar=daily
RandomizedDelaySec=5h
AccuracySec=1us

[Install]
WantedBy=timer.target
