# systemd/system/studio-background-restart.service
[Unit]
Description=restart Studio background task handler

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart studio-background.service
