# systemd/system/studio-background-restart.timer
[Timer]
OnActiveSec=48h
OnUnitActiveSec=48h

[Install]
WantedBy=timer.target
