# deploy/etc/systemd/system/studio-process-deletion-requests.timer
[Timer]
OnCalendar=daily
RandomizedDelaySec=3h
AccuracySec=1us

[Install]
WantedBy=timer.target
