# deploy/etc/systemd/system/studio-clock-tick.timer
[Timer]
OnCalendar=hourly
RandomizedDelaySec=180
AccuracySec=1us

[Install]
WantedBy=timer.target
